# infra/takedeep/supervisor/takedeep.conf
[program:takedeep-api]
environment=PROC_NAME="%(program_name)s"
command=/bin/bash -c 'set -a && source /opt/takedeep/.env && export PATH=/root/.bun/bin:$PATH && exec bun run start:api'
directory=/opt/takedeep/Deepcode-vscode/packages/takedeep-server
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[program:takedeep-dash]
environment=PROC_NAME="%(program_name)s"
command=/bin/bash -c 'set -a && source /opt/takedeep/.env && export PATH=/root/.bun/bin:$PATH && exec bun run start:dash'
directory=/opt/takedeep/Deepcode-vscode/packages/takedeep-server
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[program:takedeep-serve]
environment=PROC_NAME="%(program_name)s"
command=/bin/bash -c 'set -a && source /opt/takedeep/.env && exec /opt/takedeep/deeper serve --hostname 127.0.0.1 --port 4096'
directory=/opt/takedeep
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[program:takedeep-litellm-proxy]
environment=PROC_NAME="%(program_name)s"
command=/bin/bash -c 'set -a && source /opt/takedeep/.env && exec /opt/instance-tools/bin/cloudflared access tcp --hostname llm.qgi.dev --url 127.0.0.1:14000 --service-token-id "$CF_ACCESS_SERVICE_TOKEN_ID" --service-token-secret "$CF_ACCESS_SERVICE_TOKEN_SECRET"'
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[program:takedeep-tunnel]
environment=PROC_NAME="%(program_name)s"
command=/bin/bash -c 'set -a && source /opt/takedeep/.env && exec /opt/instance-tools/bin/cloudflared tunnel run --token "$CF_TUNNEL_SYMBOLIQ_TOKEN"'
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[program:takedeep-serve-watch]
environment=PROC_NAME="%(program_name)s"
command=/bin/bash /opt/takedeep/Deepcode-vscode/infra/takedeep/watch-serve.sh
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
