# infra/takedeep/systemd/takedeep-api.service
[Unit]
Description=TakeDeep Gateway API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/takedeep/Deepcode-vscode/packages/takedeep-server
EnvironmentFile=/opt/takedeep/.env
Environment=PATH=/root/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.bun/bin/bun run start:api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
